# ansible/playbooks/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter

[Service]
User=prometheus
ExecStart=/usr/local/bin/node_exporter {{prometheus_node_exporter_opts | default()}}

[Install]
WantedBy=default.target
